# .script/crontab
@reboot /home/ascencia/.script/on_restart.sh